# rxos/local/storage-hotplug/src/99-storage.rules
ENV{ID_FS_USAGE}=="filesystem", KERNEL!="mmcblk*", KERNEL!="loop*", RUN+="/usr/sbin/hotplug.storage"
