# rxos/local/network-config/src/99-wifi.rules
KERNEL=="wlan*", ACTION=="add", RUN+="/usr/sbin/wlan_powersave.sh"
