# rxos/rootfs/any/etc/profile.d/node_path.sh
#!/bin/sh
#
# set node path
#

export NODE_PATH=/usr/lib/node_modules
