# rxos/rootfs/any/etc/profile.d/platform-release.sh
#!/bin/sh
#
# set node path
#

. /etc/platform-release
export RXOS_SUBPLATFORM
export RXOS_VERSION
export RXOS_BUILD
export RXOS_TIMESTAMP

