# rxos/rootfs/any/etc/profile.d/user-scripts.sh
#!/bin/sh
#
# set node path
#

PATH=$PATH:/mnt/conf/bin
export PATH
