# tools/99-chip-flasher.rules
SUBSYSTEM=="usb", ATTRS{idVendor}=="1f3a", ATTRS{idProduct}=="efe8", RUN+="/path/to/autoflash.udev.sh"
